# rabbitmq/supervisor.conf
[supervisord]
nodaemon=true

[program:rabbitmq]
user=root
command=/usr/sbin/rabbitmq-server