# rabbitmq/supervisor.conf
[supervisord]
nodaemon=true

[program:rabbitmq]
command=rabbitmq-server