# rabbitmq/supervisor.conf
[supervisord]
nodaemon=true

[program:rabbitmq]
command=/opt/rabbitmq/sbin/rabbitmq-plugins enable rabbitmq_management && rabbitmq-server