# mongodb/supervisor.conf
[supervisord]
nodaemon=true

[program:mongo]
command=/opt/mongodb/bin/mongod --config /opt/mongodb/mongodb.conf